# templates/etc/systemd/system/{{name}}.service
[Unit]
Description={{description}}

[Service]
User={{name}}
ExecStart=/usr/bin/{{name}}
Type=simple
Restart=always

[Install]
WantedBy=multi-user.target
